# while installing/etc/profile.d/conda.csh
setenv CONDA_EXE "/Users/lopmmnn/Desktop/Dynamic Price Optimization Engine with Reinforcement Learning/yes/bin/conda";
setenv _CONDA_EXE "/Users/lopmmnn/Desktop/Dynamic Price Optimization Engine with Reinforcement Learning/yes/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/Users/lopmmnn/Desktop/Dynamic Price Optimization Engine with Reinforcement Learning/yes/bin/python";
setenv _CONDA_ROOT "/Users/lopmmnn/Desktop/Dynamic Price Optimization Engine with Reinforcement Learning/yes";

source "/Users/lopmmnn/Desktop/Dynamic Price Optimization Engine with Reinforcement Learning/yes/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"